# settings/service_2/supervisord.conf
# =========================
# Supervisord
# =====
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/dev/null

# =========================
# SupervisorCTL
# =====
[unix_http_server]
file=/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

# =========================
# Programs
# =====
[program:fastapi]
command=uvicorn src.fast_api.server:APP --host 0.0.0.0 --port %(ENV_SERVICE_2_HTTP_PORT)s --reload --no-access-log --proxy-headers
directory=/app
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:grpc]
command=python src/rpc/server.py
directory=/app
environment=PYTHONPATH=/app
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
